# apps/mac/JacesMac/Resources/com.jaces.mac.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jaces.mac</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/JacesMac.app/Contents/MacOS/JacesMac</string>
        <string>--headless</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>StandardOutPath</key>
    <string>/Users/USER_NAME/.jaces-mac/launchagent.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/USER_NAME/.jaces-mac/launchagent.stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/USER_NAME</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/USER_NAME</string>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>